# punchtime_web.service
[Unit]
Description=Punchtime

[Service]
ExecStart=/usr/local/bin/punchtime serve
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
