# punchtime_recorder.service
[Unit]
Description=Punchtime

[Service]
ExecStart=/usr/local/bin/punchtime iterate
