# punchtime_web.service

[Unit]
Description=Punchtime

[Service]
ExecStart=/usr/local/bin/punchtime serve
Restart=always
